# stage2/04-comitup/files/wifi-on.service
[Unit]
Description=Turn wifi on, regardless of regulatory domain
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/nmcli radio wifi on

[Install]
WantedBy=NetworkManager.service
